# debian/baruwa.cron.daily
#!/bin/sh
#
# Baruwa
#
# send quarantine reports
/usr/share/pyshared/baruwa/manage.py sendquarantinereports
# clean quarantine 
/usr/share/pyshared/baruwa/manage.py cleanquarantine
# clean up the DB
/usr/share/pyshared/baruwa/manage.py dbclean
# update sa rule definitions
#/usr/share/pyshared/baruwa/manage.py update_sarules
# clean up stale sessions
/usr/share/pyshared/baruwa/manage.py cleanup
# update geoip database
/usr/share/doc/libgeoip1/examples/geolitecountryv4.sh
